# backends/agentmux/agentmux-update.timer.tmpl
[Unit]
Description=Periodic maintenance for agentmux instance @@INSTANCE_NAME@@

[Timer]
OnCalendar=@@ON_CALENDAR@@
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
