# backends/agentmux/agentmux.service.tmpl
[Unit]
Description=Persistent agentmux instance @@INSTANCE_NAME@@ (@@AGENT@@ + @@PROVIDER@@)
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=@@RUN_USER@@
EnvironmentFile=@@ENV_FILE@@
ExecStart=@@REPO_DIR@@/rc-start.sh
ExecStop=/usr/bin/tmux -L agentmux-@@INSTANCE_NAME@@ kill-session -t @@TMUX_SESSION_NAME@@
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
