# backends/claude-code/agentmux-claude-code-update.timer.tmpl
[Unit]
Description=Periodic Claude Code update for @@INSTANCE_NAME@@ / @@SESSION_NAME@@

[Timer]
OnCalendar=@@ON_CALENDAR@@
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
