# backends/claude-code/agentmux-claude-code-update.service.tmpl
[Unit]
Description=Update Claude Code CLI and restart the @@SESSION_NAME@@ session
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=@@ENV_FILE@@
ExecStart=@@REPO_DIR@@/rc-update.sh
