# backends/claude-code/agentmux-claude-code.service.tmpl
[Unit]
Description=Persistent agentmux Claude Code session (@@INSTANCE_NAME@@ / @@SESSION_NAME@@)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=@@RUN_USER@@
EnvironmentFile=@@ENV_FILE@@
ExecStart=/bin/bash @@REPO_DIR@@/rc-start.sh
ExecStop=/bin/bash -c 'tmux -L agentmux-@@INSTANCE_NAME@@ kill-session -t @@SESSION_NAME@@'
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
