# backends/claude-code/claude-code.update.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>@@UPDATE_LABEL@@</string>

    <key>ProgramArguments</key>
    <array>
        <string>@@REPO_DIR@@/rc-update.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>@@UPDATE_HOUR@@</integer>
        <key>Minute</key>
        <integer>@@UPDATE_MINUTE@@</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>AGENTMUX_INSTANCE_NAME</key>
        <string>@@INSTANCE_NAME@@</string>
        <key>AGENTMUX_SESSION_NAME</key>
        <string>@@TMUX_SESSION_NAME@@</string>
        <key>AGENTMUX_TMUX_SESSION_NAME</key>
        <string>@@TMUX_SESSION_NAME@@</string>
        <key>AGENTMUX_DISPLAY_NAME</key>
        <string>@@DISPLAY_NAME@@</string>
        <key>AGENTMUX_WORKDIR</key>
        <string>@@WORKDIR@@</string>
        <key>AGENTMUX_RESUME</key>
        <string>@@RESUME@@</string>
        <key>HOME</key>
        <string>@@HOME@@</string>
        <key>PATH</key>
        <string>@@PATH@@</string>
    </dict>

    <key>StandardOutPath</key>
    <string>@@LOG_DIR@@/@@INSTANCE_NAME@@-update.log</string>
    <key>StandardErrorPath</key>
    <string>@@LOG_DIR@@/@@INSTANCE_NAME@@-update.err.log</string>
</dict>
</plist>
